# .config/theme-setter/set-theme.service
[Unit]
Description=Global Theme Setter Service

[Service]
ExecStart=set-theme $(get-theme) $(get-theme --mode)

[Install]
WantedBy=multi-user.target